# dist/lib/systemd/system/fxaudio.service
[Unit]
Description=fxAudio service

[Service]
User=fxaudio
Group=audio
ExecStart=/usr/bin/fxaudio start --config=/etc/fxaudio/config.yaml
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
