# dist/lib/systemd/system/fxaudio.service
[Unit]
Description=fxAudio service
Wants=sound.target

[Service]
User=admin
Group=admin
ExecStart=/usr/bin/fxaudio start --config=/etc/fxaudio/config.yaml
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
