# dist/lib/systemd/system/fxaudio.service
[Unit]
Description=fxAudio service

[Service]
User=root
Group=root
ExecStart=/usr/bin/fxaudio start --config=/etc/fxaudio/config.yaml

[Install]
WantedBy=multi-user.target