# dist/DEBIAN/fxaudio.service.in
[Unit]
Description=fxAudio service
Wants=sound.target

[Service]
Environment="XDG_RUNTIME_DIR=/run/user/@SERVICE_UID@"
User=@SERVICE_USER@
Group=@SERVICE_USER@
ExecStart=/usr/bin/fxaudio start --config=/etc/fxaudio/config.yaml
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target

