# dist/DEBIAN/postinst
#!/bin/sh
systemctl enable fxaudio.service
systemctl start fxaudio.service