# dist/DEBIAN/postinst
#!/bin/sh
set -e

# Get the logged in user as the default user
current_user=$(logname)

# Prompt the user
printf "Enter the username that should run the fxaudio service [default: %s]" "$current_user"
read service_user

# Use current user if nothing entered
if [ -z "$service_user" ]; then
    service_user="$current_user"
fi

echo "Installing service to run as user: $service_user"

# Get the uid of the user
service_uid=$(id -u "$service_user")
echo "User id for $service_user is $service_uid"

# Ensure data directory exists and is writable by service user
mkdir -p /var/lib/fxaudio
chown "$service_user:$service_user" /var/lib/fxaudio
chmod 755 /var/lib/fxaudio

echo "Created data folder /var/lib/fxaudio"

# Path to your template and destination
TEMPLATE_PATH="/usr/lib/fxaudio/fxaudio.service.in"
OUTPUT_PATH="/lib/systemd/system/fxaudio.service"

echo "Creating systemd unit file $OUTPUT_PATH"

# Do the replacements
sed \
  -e "s|@SERVICE_UID@|$service_uid|g" \
  -e "s|@SERVICE_USER@|$service_user|g" \
  "$TEMPLATE_PATH" > "$OUTPUT_PATH"

# Reload, enable, and start the service
chmod 644 /lib/systemd/system/fxaudio.service
systemctl daemon-reload
systemctl enable fxaudio.service
systemctl start fxaudio.service


