# dist/DEBIAN/prerm
#!/bin/sh
systemctl stop fxaudio.service
systemctl disable fxaudio.service
sudo userdel fxaudio